# roles/ktz-cartman/files/etc/systemd/system/docker.service.d/zfs-override.conf
[Unit]
Requires=zfs.target
After=zfs.target