# roles/monitoring-exporters/templates/node-intel-gpu-textfile.service.j2
[Unit]
Description=Write Intel GPU metrics for Prometheus node_exporter
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/node-intel-gpu-textfile
