# roles/monitoring-exporters/templates/node-smart-textfile.service.j2
[Unit]
Description=Write SMART metrics for Prometheus node_exporter
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/node-smart-textfile
