# legacy/roles/systemd-timer-plex/templates/plex-restart.timer.j2
[Unit]
Description=Timer to restart Plex Docker container
Requires=plex-restart.service

[Timer]
OnCalendar=*-*-* {{ plex_restart_time }}
Persistent=true

[Install]
WantedBy=timers.target